# chef-repo/cookbooks/kubernetes/files/default/services/kube-apiserver.service
[Unit]
Description=Kubernetes API Server
After=etcd.service
After=docker.service
Wants=etcd.service
Wants=docker.service

[Service]

