# chef-repo/cookbooks/pod_manager/files/default/dockerfiles/apps/zookeeper/supervisor.conf
[supervisord]
nodaemon=false

[program:chef-client]
command=/usr/bin/chef-client -r "role[zookeeper]"
